# scripts/systemd/exim4@1.service.d/overrides.conf
[Unit]
Description=Incoming MTA stage for SpamTagger
After=exim4@2.service exim4@4.service mariadb@replica.service
Requires=exim4@2.service exim4@4.service mariadb@replica.service
PartOf=spamtagger.target
BindsTo=greylistd.service 

[Service]
User=spamtagger
Group=spamtagger
ExecStart=/opt/exim4/bin/exim -C /var/spamtagger/spool/tmp/exim/exim_stage%I.conf -bd -q5m

[Install]
WantedBy=spamtagger.target
