# scripts/systemd/mailcleaner-stopped.target
[Unit]
Description=SpamTagger - SMTP Antivirus/Antispam Gateway
Requires=sshd.service
Conflicts=spamtagger.target
AllowIsolate=yes
